# root/etc/uci-defaults/40_luci-ss-domain-rule
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ss-domain-rule[-1]
	add ucitrack ss-domain-rule
	set ucitrack.@ss-domain-rule[-1].init=ss-domain-rule
	commit ucitrack
EOF
